# tajvidbot.conf
[program:tajvidbot]
directory=/home/ubuntu/Tajvid_tekshirish
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/Tajvid_tekshirish/logfile_err.log
stdout_logfile=/home/ubuntu/Tajvid_tekshirish/logfile.log